# honeypot.service
[Unit]
Description=honeypot daemon

[Service]
Type=simple
ExecStart=/bin/bash -c /opt/honeypot/start.sh

[Install]
WantedBy=multi-user.target